# deploy/pr_dashboard-gunicorn.service
[Unit]
Description=Gunicorn for pr_dashboard Django backend
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/pr_dashborad/backend
EnvironmentFile=/var/www/pr_dashborad/backend/.env
ExecStart=/var/www/pr_dashborad/backend/.venv/bin/gunicorn -c /var/www/pr_dashborad/backend/gunicorn.conf.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
